# scripts/upstart/statique.conf
description "statique"
start on stopped rc RUNLEVEL=[2345]
respawn
chdir /home/ubuntu/statique
exec /home/ubuntu/statique/statique.bin -p 80
